# systemd/moosefs-cgiserv.service.in
[Unit]
Description=MooseFS CGI server
Wants=network-online.target
After=network.target network-online.target

[Service]
Type=forking
ExecStart=@SBIN_PATH@/mfscgiserv start
ExecStop=@SBIN_PATH@/mfscgiserv stop
PIDFile=@DATA_PATH@/.mfscgiserv.lock
Restart=on-abnormal

[Install]
WantedBy=multi-user.target
