# freebsd/files/files.master/mfsmaster.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: mfsmaster
# REQUIRE: DAEMON
# KEYWORD: nostart shutdown
# BEFORE: mountlate
#
# Add the following lines to /etc/rc.conf to enable mfsmaster:
#
# mfsmaster_enable="YES"
#

. /etc/rc.subr

name=mfsmaster
rcvar=mfsmaster_enable

config_file=%%PREFIX%%/etc/mfs/${name}.cfg

required_files=${config_file}

command=%%PREFIX%%/sbin/${name}
pid_file=%%MFS_WORKDIR%%/mfs/.${name}.lock
command_args="-c $config_file"

load_rc_config $name

# set defaults
: ${mfsmaster_enable="NO"}

run_rc_command "$1"
